# config/systemd/chefoftheday-holidays.timer
[Unit]
Description=ChefOfTheDay Feiertage Timer (jährlich 1. Juli)

[Timer]
OnCalendar=*-07-01 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
